# clusterware-fshelper/etc/profile.d/80-fshelper.csh
################################################################################
##
## Alces Clusterware - Shell configuration
## Copyright (c) 2017 Alces Software Ltd
##
################################################################################
if ( -x "$cw_ROOT/libexec/share/setup-fslinks" ) then
    $cw_ROOT/libexec/share/setup-fslinks
endif
